# scripts/rpm/after-install
#!/bin/bash

cp -u /etc/kiwiirc/config.conf.example /etc/kiwiirc/config.conf
chkconfig --add webircgateway
chkconfig webircgateway on
/etc/init.d/webircgateway start

