# scripts/rpm/after-install
#!/bin/bash

chkconfig --add webircgateway
chkconfig webircgateway on
/etc/init.d/webircgateway start

