# scripts/rpm/before-remove
#!/bin/bash

/etc/init.d/webircgateway stop
chkconfig webircgateway off
chkconfig --del webircgateway

