# scripts/rpm/before-remove
#!/bin/bash

/etc/init.d/kiwiirc stop
chkconfig kiwiirc off
chkconfig --del kiwiirc

