# scripts/deb/after-install
#!/bin/bash

update-rc.d webircgateway defaults
/etc/init.d/webircgateway start
