# scripts/deb/after-install
#!/bin/bash

adduser --system --group --home /var/lib/kiwiirc kiwiirc

cp -u /etc/kiwiirc/config.conf.example /etc/kiwiirc/config.conf
update-rc.d kiwiirc defaults
/etc/init.d/kiwiirc start
