# scripts/deb/before-remove
#!/bin/bash

/etc/init.d/kiwiirc stop
update-rc.d -f kiwiirc remove
