# scripts/deb/before-remove
#!/bin/bash

/etc/init.d/webircgateway stop
update-rc.d -f webircgateway remove
